# bareos/core/platforms/univention/univention-bareos.cron
# restart the bareos director before running the daily backup jobs
30 20 * * * root /usr/share/univention-bareos/bareos-dir-restart-on-change
