# core/platforms/osx/files/org.bareos.bareos-fd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>org.bareos.bareos-fd</string>

        <key>ProgramArguments</key>
        <array>
                <string>@sbindir@/bareos-fd</string>
                <string>-f</string>
        </array>

        <key>StandardOutPath</key>
        <string>/var/run/bareos-fd.log</string>

        <key>StandardErrorPath</key>
        <string>/var/run/bareos.log</string>

        <key>RunAtLoad</key>
        <true/>
</dict>
</plist>
