# core/platforms/gentoo/bareos-fd.in
#!/sbin/runscript
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos File daemon for the Gentoo release.
#
# chkconfig: 2345 90 9
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#

start() {
       ebegin "Starting the Bareos File daemon"
       start-stop-daemon --start --quiet --exec @sbindir@/bareos-fd -- $2 -c @confdir@/bareos-fd.conf
       eend $?
}

stop() {
       ebegin "Stopping the Bareos File daemon"
       start-stop-daemon --stop --quiet --exec @sbindir@/bareos-fd
       eend $?
}

restart() {
       stop
       sleep 5
       start
}
