# core/platforms/gentoo/bareos-init.in
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# Modified for user/group information
# 24 Oct 2004 D. Scott Barninger <barninger at fairfieldcomputers dot com

depend() {
	need net
	use @DEFAULT_DB_TYPE@
}

opts="test"

DIR_USER=@dir_user@
DIR_GROUP=@dir_group@
SD_USER=@sd_user@
SD_GROUP=@sd_group@
FD_USER=@fd_user@
FD_GROUP=@fd_group@

# A client would only need to run the file daemon, so if thats the
# case for this machine just comment the other two out below here...
start() {
	ebegin "Starting bareos storage daemon"
	OPTIONS=''
	if [ "${SD_USER}" != '' ]; then
		OPTIONS="${OPTIONS} -u ${SD_USER}"
	fi

	if [ "${SD_GROUP}" != '' ]; then
		OPTIONS="${OPTIONS} -g ${SD_GROUP}"
	fi
	start-stop-daemon --start --quiet --exec @sbindir@/bareos-sd \
		-- ${OPTIONS} -c @confdir@/bareos-sd.conf
	result=$?

	ebegin "Starting bareos file daemon"
	OPTIONS=''
	if [ "${FD_USER}" != '' ]; then
		OPTIONS="${OPTIONS} -u ${FD_USER}"
	fi

	if [ "${FD_GROUP}" != '' ]; then
		OPTIONS="${OPTIONS} -g ${FD_GROUP}"
	fi
	start-stop-daemon --start --quiet --exec @sbindir@/bareos-fd \
		-- ${OPTIONS} -c @confdir@/bareos-fd.conf
	result=$(( $result + $? ))
	sleep 2

	ebegin "Starting bareos director"
	OPTIONS=''
	if [ "${DIR_USER}" != '' ]; then
		OPTIONS="${OPTIONS} -u ${DIR_USER}"
	fi

	if [ "${DIR_GROUP}" != '' ]; then
		OPTIONS="${OPTIONS} -g ${DIR_GROUP}"
	fi
	start-stop-daemon --start --quiet --exec @sbindir@/bareos-dir \
		-- ${OPTIONS} -c @confdir@/bareos-dir.conf
	result=$(( $result + $? ))

	eend $result
}

# The default port numbers appear in the pidfile paths below.
# Might wanna programmatically grab them from the config files
# or something...
stop() {
	ebegin "Stopping bareos file daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/bareos-fd.9102.pid
	result=$?
	ebegin "Stopping bareos storage daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/bareos-sd.9103.pid
	result=$(( $result + $? ))
	ebegin "Stopping bareos director"
	start-stop-daemon --stop --quiet --pidfile /var/run/bareos-dir.9101.pid
	result=$(( $result + $? ))
	eend $result
}

test() {
	ebegin "Testing bareos storage daemon configuration"
	@sbindir@/bareos-sd -t -c @confdir@/bareos-sd.conf
	result=$?
	ebegin "Testing bareos file daemon configuration"
	@sbindir@/bareos-fd -t -c @confdir@/bareos-fd.conf
	result=$(( $result + $? ))
	ebegin "Testing bareos director configuration"
	@sbindir@/bareos-dir -t -c @confdir@/bareos-dir.conf
	result=$(( $result + $? ))
	eend $result
}
