# platforms/gentoo/bacula-dir.in
#!/sbin/runscript
#
# bacula       This shell script takes care of starting and stopping
#	       the bacula Director daemon for the Gentoo release
#
# chkconfig: 2345 90 9
# description: The Leading Open Source Backup Solution.
#
#  For Bacula release @VERSION@ (@DATE@) -- @DISTNAME@
#

start() {
       ebegin "Starting the Bacula Director"
       start-stop-daemon --start --quiet --exec @sbindir@/bacula-dir -- $2 -c @sysconfdir@/bacula-dir.conf
       eend $?
}

stop() {
       ebegin "Stopping the Director daemon"
       start-stop-daemon --stop --quiet --exec @sbindir@/bacula-dir
       eend $?
}

restart() {
       stop
       sleep 5
       start
}
