# platforms/gentoo/bareos-sd.in
#!/sbin/runscript
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos Storage daemon for the Gentoo release.
#
# chkconfig: 2345 90 9
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#

start() {
       ebegin "Starting the Bareos Storage daemon"
       start-stop-daemon --start --quiet --exec  @sbindir@/bareos-sd -- $2 -c @confdir@/bareos-sd.conf
       eend $?
}

stop() {
       ebegin "Stopping the Bareos Storage daemon"
       start-stop-daemon --stop --quiet --exec @sbindir@/bareos-sd
       eend $?
}

restart() {
       stop
       sleep 5
       start
}
