# platforms/univention/univention-bareos.cron
30 0 * * * /usr/share/univention-bareos/restart_director
