# platforms/univention/univention-bareos.cron
# restart the bareos director before running the daily backup jobs
30 0 * * * root /usr/share/univention-bareos/restart_director
