# platforms/osx/files/org.bacula.bacula-fd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>org.bacula.bacula-fd</string>
        <key>ProgramArguments</key>
        <array>
                <string>@PREFIX@/sbin/bacula-fd</string>
                <string>-f</string>
                <string>-c</string>
                <string>@FD_CONF@</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
</dict>
</plist>
