# core/platforms/systemd/bareos-sd.service.in
#
# Copyright (C) 2013-2023 Bareos GmbH & Co. KG
# Copyright (c) 2011 Free Software Foundation Europe e.V., Author: Bruno Friedmann
#
[Unit]
Description=Bareos Storage Daemon service
Documentation=man:bareos-sd(8)
Requires=network.target nss-lookup.target time-sync.target
After=network-online.target nss-lookup.target remote-fs.target time-sync.target
# Check if working dir exist
ConditionPathExists=@working_dir@
# Limit the number of restart per day to 10
StartLimitIntervalSec=1d
StartLimitBurst=10

[Service]
Type=simple
User=@sd_user@
Group=@sd_group@
WorkingDirectory=@working_dir@
ExecStart=@sbindir@/bareos-sd -f
SuccessExitStatus=0 15
# Increase the the maximum number of open file descriptors
LimitNOFILE=8192:524288
# Restart on failure, wait 30 seconds
Restart=on-failure
RestartSec=30
# Don't restart on wrong exec arguments and configuration errors.
RestartPreventExitStatus=109 1


[Install]
Alias=bareos-storage.service
WantedBy=multi-user.target
