# core/platforms/systemd/bareos-dir.service.in
#
# Copyright (C) 2013-2021 Bareos GmbH & Co. KG
# Copyright (c) 2011 Free Software Foundation Europe e.V., Author: Bruno Friedmann
#
[Unit]
Description=Bareos Director Daemon service
Documentation=man:bareos-dir(8)
Requires=network.target nss-lookup.target time-sync.target
After=network-online.target nss-lookup.target remote-fs.target time-sync.target postgresql.service mysql.service mariadb.service
# Dependency about the database
# We let administrators decide if they need it (if local db instance)
# Wants=@DEFAULT_DB_TYPE@.service
# Check if working dir exist
ConditionPathExists=@working_dir@
# Limit the number of restart per day to 10
# wait 3 minutes between each restart
RestartSec=180
StartLimitIntervalSec=1d
StartLimitBurst=10

[Service]
Type=simple
User=@dir_user@
Group=@dir_group@
WorkingDirectory=@working_dir@
ExecStart=@sbindir@/bareos-dir -f
SuccessExitStatus=0 15
ExecReload=@sbindir@/bareos-dir -t -f
ExecReload=/bin/kill -HUP $MAINPID
# Increase the maximum number of open file descriptors
LimitNOFILE=8192:524288

[Install]
Alias=bareos-director.service
WantedBy=multi-user.target
