# platforms/systemd/bareos-sd.service.in
# This file is part of package Bareos Storage Daemon
#
# Copyright (c) 2011 Free Software Foundation Europe e.V.
# for Bareos Community
# Author: Bruno Friedmann
# Description:
#    Used to start the bareos storage daemon service (bareos-sd)
#    will be installed as /lib/systemd/system/bareos-sd.service
#    enable : systemctl enable bareos-sd.service
#    start : systemctl start bareos-sd.service
#
# Bareos Storage Daemon service
#
[Unit]
Description=Bareos Storage Daemon service
Documentation=man:bareos-sd(8)
Requires=nss-lookup.target network.target remote-fs.target time-sync.target
After=nss-lookup.target network.target remote-fs.target time-sync.target
# Wants=
# Before=
# Conflicts=

[Service]
Type=forking
User=@sd_user@
Group=@sd_group@
WorkingDirectory=@working_dir@
PIDFile=@piddir@/bareos-sd.@sd_port@.pid
# EnvironmentFile=-/etc/sysconfig/bareos-sd
StandardOutput=syslog
ExecStart=@sbindir@/bareos-sd
# enable this for scsicrypto-sd
# CapabilityBoundingSet=cap_sys_rawio+ep
SuccessExitStatus=0 15
Restart=on-failure

[Install]
Alias=bareos-storage.service
WantedBy=multi-user.target
