# platforms/systemd/bareos-sd.socket
[Unit]
Conflicts=bareos-sd.service

[Socket]
ListenStream=9103
Accept=yes

[Install]
WantedBy=sockets.target
