# platforms/systemd/bareos-fd.socket
[Unit]
Conflicts=bareos-fd.service

[Socket]
ListenStream=9102
Accept=yes

[Install]
WantedBy=sockets.target
