# core/platforms/slackware/rc.bareos-sd.in
#! /bin/sh
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos Storage daemon.
#
# chkconfig: 2345 90 99
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#

# Source function library
. /etc/rc.d/functions.bareos

case "$1" in
    start)
	[ -x ${BACSDBIN}/bareos-sd ] && {
	    sleep 2
	    echo -n "Starting the Storage daemon: "
	    OPTIONS=''
	    if [ "${SD_USER}" != '' ]; then
	       OPTIONS="${OPTIONS} -u ${SD_USER}"
	    fi

	    if [ "${SD_GROUP}" != '' ]; then
	       OPTIONS="${OPTIONS} -g ${SD_GROUP}"
	    fi

	    ${BACSDBIN}/bareos-sd $2 ${OPTIONS} -v -c ${BACSDCFG}/bareos-sd.conf
	    echo Done.
	}
	;;
    stop)
	[ -x ${BACSDBIN}/bareos-sd ] && {
	    echo -n "Stopping the Storage daemon: "
	    killproc ${BACSDBIN}/bareos-sd ${SD_PORT}
	    echo Done.
	}
	;;
    restart)
	$0 stop
	sleep 5
	$0 start
	;;
    status)
	[ -x ${BACSDBIN}/bareos-sd ] && status ${BACSDBIN}/bareos-sd ${SD_PORT}
	;;
    *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
	;;
esac
exit 0
