# core/platforms/slackware/rc.bareos-dir.in
#! /bin/sh
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos Director daemon
#
# chkconfig: 2345 92 99
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#

# Source function library
. /etc/rc.d/functions.bareos

case "$1" in
    start)
	[ -x ${BACDIRBIN}/bareos-dir ] && {
	    sleep 2
	    echo -n "Starting the Director daemon: "
	    OPTIONS=''
	    if [ "${DIR_USER}" != '' ]; then
	       OPTIONS="${OPTIONS} -u ${DIR_USER}"
	    fi

	    if [ "${DIR_GROUP}" != '' ]; then
	       OPTIONS="${OPTIONS} -g ${DIR_GROUP}"
	    fi

	    ${BACDIRBIN}/bareos-dir $2 ${OPTIONS} -v -c ${BACDIRCFG}/bareos-dir.conf
	    echo Done.
	}
	;;
    stop)
	[ -x ${BACDIRBIN}/bareos-dir ] && {
	    echo -n "Stopping the Director daemon: "
	    killproc ${BACDIRBIN}/bareos-dir ${DIR_PORT}
	    echo Done.
	}
	;;
    restart)
	$0 stop
	sleep 5
	$0 start
	;;
    status)
	[ -x ${BACDIRBIN}/bareos-dir ] && status ${BACDIRBIN}/bareos-dir ${DIR_PORT}
	;;
    *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
	;;
esac
exit 0
