# bareos/core/platforms/systemd/bareos-fd.service.in
#
# Copyright (C) 2013-2021 Bareos GmbH & Co. KG
# Copyright (c) 2011 Free Software Foundation Europe e.V., Author: Bruno Friedmann
#
[Unit]
Description=Bareos File Daemon service
Documentation=man:bareos-fd(8)
Requires=network.target nss-lookup.target time-sync.target
After=network-online.target nss-lookup.target remote-fs.target time-sync.target

[Service]
Type=simple
User=@fd_user@
Group=@fd_group@
WorkingDirectory=@working_dir@
ExecStart=@sbindir@/bareos-fd -f
SuccessExitStatus=0 15
Restart=on-failure
# IOSchedulingClass=idle

[Install]
Alias=bareos-filedaemon.service
WantedBy=multi-user.target
