# platforms/mandrake/bacula-sd.in
#! /bin/sh
#
# bacula       This shell script takes care of starting and stopping
#	       the bacula Storage daemon.
#
# chkconfig: 2345 90 9
# description: The Leading Open Source Backup Solution.
#
#  For Bacula release @VERSION@ (@DATE@) -- @DISTNAME@
#

# Source function library
. /etc/rc.d/init.d/functions

case "$1" in
    start)
       echo -n "Starting the Bacula Storage daemon: "
       daemon @sbindir@/bacula-sd $2 -c @sysconfdir@/bacula-sd.conf
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && touch @subsysdir@/bacula-sd
       ;;
    stop)
       echo -n "Stopping the Bacula Storage daemon: "
       killproc @sbindir@/bacula-sd
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && rm -f @subsysdir@/bacula-sd
       ;;
    restart)
       $0 stop
       sleep 5
       $0 start
       ;;
    status)
       status @sbindir@/bacula-sd
       ;;
    *)
       echo "Usage: $0 {start|stop|restart|status}"
       exit 1
       ;;
esac
exit 0
