# core/platforms/freebsd/bareos-freebsd/bareos.com-common/pkg-plist.storage
@dir lib/bareos/defaultconfigs/bareos-sd.d/
@dir lib/bareos/defaultconfigs/bareos-sd.d/autochanger/
@dir lib/bareos/defaultconfigs/bareos-sd.d/device/
@dir lib/bareos/defaultconfigs/bareos-sd.d/director/
@dir lib/bareos/defaultconfigs/bareos-sd.d/ndmp/
@dir lib/bareos/defaultconfigs/bareos-sd.d/messages/
@dir lib/bareos/defaultconfigs/bareos-sd.d/storage/
lib/bareos/defaultconfigs/bareos-sd.d/autochanger/FileStorage.conf.example
lib/bareos/defaultconfigs/bareos-sd.d/device/FileStorage.conf.example
lib/bareos/defaultconfigs/bareos-sd.d/device/FileStorage.conf
lib/bareos/defaultconfigs/bareos-sd.d/director/bareos-dir.conf
lib/bareos/defaultconfigs/bareos-sd.d/director/bareos-mon.conf
lib/bareos/defaultconfigs/bareos-sd.d/messages/Standard.conf
lib/bareos/defaultconfigs/bareos-sd.d/storage/bareos-sd.conf
sbin/bareos-sd
lib/bareos/scripts/disk-changer
lib/bareos/plugins/autoxflate-sd.so
lib/libbareossd-file.so
share/man/man8/bareos-sd.8.gz
@dir(bareos,bareos,0775) /var/lib/bareos/storage
etc/rc.d/bareos-sd
