# core/platforms/freebsd/bareos-sd.in
#! /bin/sh
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos Storage daemon.
#
# chkconfig: 2345 90 9
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#
. /etc/rc.subr

name=bareossd
rcvar=bareossd_enable

command="@sbindir@/bareos-sd"

load_rc_config $name
run_rc_command "$1"

