# core/platforms/freebsd/bareos-fd.in
#!/bin/sh
#
# bareos       This shell script takes care of starting and stopping
#	       the bareos File daemon.
#
# chkconfig: 2345 90 9
# description: Backup Archiving REcovery Open Sourced.
#
#  For Bareos release @VERSION@ (@DATE@) -- @DISTNAME@
#
. /etc/rc.subr

name=bareosfd
rcvar=bareosfd_enable

command="@sbindir@/bareos-fd"

load_rc_config $name
run_rc_command "$1"

