# debian/bareos-storage.postinst.in
#!/bin/sh
# postinst script for bareos
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

daemon_user=bareos
daemon_group=bareos

#director_daemon_user=$daemon_user
storage_daemon_user=$daemon_user
#file_daemon_user=root
storage_daemon_group=$daemon_group

DEFCONFIGDIR="@configtemplatedir@/bareos-sd.d/"
CONFIGDIR="/etc/bareos/bareos-sd.d/"



permissions()
{
  chown ${storage_daemon_user}:${daemon_group} ${CONFIGDIR}
  chmod 750 ${CONFIGDIR}
  chown ${storage_daemon_user}:${daemon_group} /var/lib/bareos/storage/
  chmod 750 /var/lib/bareos/storage/
}

case "$1" in
    configure)
        if ! /usr/lib/bareos/scripts/bareos-config has_config "bareos-sd"; then
            /usr/lib/bareos/scripts/bareos-config deploy_config "$DEFCONFIGDIR" "$CONFIGDIR"
            /usr/lib/bareos/scripts/bareos-config initialize_local_hostname
            /usr/lib/bareos/scripts/bareos-config initialize_passwords
        fi
        permissions
        /usr/lib/bareos/scripts/bareos-config setup_sd_user
        # on Univention distributions the ucr command allows to open the firewall
        if [ -x "/etc/init.d/univention-firewall" ] && which ucr >/dev/null 2>&1; then
            ucr set \
                security/packetfilter/package/bareos-storage/tcp/@sd_port@/all="ACCEPT" \
                security/packetfilter/package/bareos-storage/tcp/@sd_port@/all/en="bareos-sd"
            #invoke-rc.d univention-firewall restart
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
