# debian/bareos-storage.service.in
#
# Copyright (C) 2013-2023 Bareos GmbH & Co. KG
# Copyright (c) 2011 Free Software Foundation Europe e.V., Author: Bruno Friedmann
#
[Unit]
Description=Bareos Storage Daemon service
Documentation=man:bareos-sd(8)
Requires=network.target nss-lookup.target time-sync.target
After=network-online.target nss-lookup.target remote-fs.target time-sync.target
# Check if working dir exist
ConditionPathExists=@working_dir@
# Limit the number of restart per day to 10
# wait 3 minutes between each restart
RestartSec=180
StartLimitIntervalSec=1d
StartLimitBurst=10

[Service]
Type=simple
User=@sd_user@
Group=@sd_group@
WorkingDirectory=@working_dir@
ExecStart=@sbindir@/bareos-sd -f
SuccessExitStatus=0 15
# Increase the maximum number of open file descriptors```
LimitNOFILE=8192:524288

[Install]
Alias=bareos-sd.service
WantedBy=multi-user.target
